# daemons/com.clavier.ray-dalio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.clavier.ray-dalio</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/doppler</string>
        <string>run</string>
        <string>--project</string>
        <string>clavier</string>
        <string>--config</string>
        <string>prd</string>
        <string>--silent</string>
        <string>--</string>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/clavier/Library/Mobile Documents/com~apple~CloudDocs/0/scripts/tools/ray-dalio-runner.mjs</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>30</integer></dict>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/Users/clavier/Library/Logs/ray-dalio.log</string>
    <key>StandardErrorPath</key><string>/Users/clavier/Library/Logs/ray-dalio_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key><string>/Users/clavier</string>
    </dict>
</dict>
</plist>
