# launchagents/com.clavier.scripts-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clavier.scripts-watcher</string>

    <!-- iCloud/scripts/ 폴더 변경 감지 시 install-scripts.sh 자동 실행 -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/clavier/Library/Mobile Documents/com~apple~CloudDocs/0/scripts</string>
    </array>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/clavier/Library/Mobile Documents/com~apple~CloudDocs/0/scripts/install-scripts.sh</string>
    </array>

    <!-- 부팅 시 한 번 실행해서 ~/bin 초기 동기화 -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/clavier/Library/Logs/scripts-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/clavier/Library/Logs/scripts-watcher-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/clavier</string>
    </dict>
</dict>
</plist>
